# systemd/moondev-monitor.timer
[Unit]
Description=Moon Dev Monitor weekly timer

[Timer]
OnCalendar=Sun *-*-* 10:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
